# installscripts/dockerfiles/jenkins/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=service ssh start

[program:jenkins]
command=service jenkins start
